# docker_config_files/supervisor_nginx_proto_app.conf
[supervisord]
nodaemon=true

[program:gunicorn]
directory=/home/app/proto_app
command=gunicorn proto_app.wsgi:application -b localhost:8000


[program:nginx]
command=/usr/sbin/nginx